# Utils/crontab
0 */2 * * 1-5 root /Collectors/start.sh > /proc/1/fd/1 2>/proc/1/fd/2
